# src/templates/services/robonet-monitor.service
[Unit]
Description=Robotnet Monitor
After=network.target

[Service]
Type=simple
ExecStart=${EXECUTABLE_PATH} run-service
Environment=RUST_LOG=info
Environment=NOTIFICATION_EMAIL=${NOTIFICATION_EMAIL}
Environment=SMTP_SERVER=${SMTP_SERVER}
Environment=SMTP_USER=${SMTP_USER}
Environment=SMTP_PASSWORD=${SMTP_PASSWORD}
Environment=CHECK_INTERVAL_SECS=${CHECK_INTERVAL_SECS}
Environment=MAX_RETRIES=${MAX_RETRIES}
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target 